# docker/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0
pidfile=/run/supervisord.pid

# ---- PHP-FPM ----------------------------------------------------------------
[program:php-fpm]
command=php-fpm -F
autostart=true
autorestart=true
priority=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# ---- Nginx ------------------------------------------------------------------
[program:nginx]
command=nginx -g 'daemon off;'
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# ---- Queue worker (QUEUE_CONNECTION=database) -------------------------------
[program:queue]
command=php /var/www/html/artisan queue:work --sleep=3 --tries=3 --max-time=3600 --no-interaction
user=www-data
autostart=true
autorestart=true
priority=20
stopwaitsecs=3600
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# ---- Scheduler (runs due tasks every minute) --------------------------------
[program:scheduler]
command=php /var/www/html/artisan schedule:work --no-interaction
user=www-data
autostart=true
autorestart=true
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
